# com.user.music-format-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.music-format-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Adjust this path if you install the binary elsewhere -->
        <string>~/Workspace/coreaudiohacks/music-format-daemon</string>
        <string>D10s</string>
    </array>

    <!-- Start automatically at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logs: visible in Console.app under "music-format-daemon" -->
    <key>StandardOutPath</key>
    <string>/tmp/music-format-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/music-format-daemon.err</string>
</dict>
</plist>
